# .files/lowspc/main/peditxi
#!/bin/sh /etc/rc.common

start(){

sh /root/f2.sh

}

reload(){

sh /root/opt.sh

}
